# scripts/com.news-ingestion.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.news-ingestion.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/cyc/Desktop/相关文档/00-项目/量化交易/news-ingestion-v1/scripts/collect_daily.sh</string>
  </array>

  <!-- 周一至周五 16:00（A股收盘后）各跑一次 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <!-- 若到点时机器在睡眠/关机，唤醒后补跑一次 -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/cyc/Desktop/相关文档/00-项目/量化交易/news-ingestion-v1/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/cyc/Desktop/相关文档/00-项目/量化交易/news-ingestion-v1/logs/launchd.err.log</string>
</dict>
</plist>
